# scripts/confs/gsc-api.socket
[Unit]
Description=gsc-api gunicorn socket

[Socket]
ListenStream=/run/gsc-api.sock

[Install]
WantedBy=sockets.target

